# crontab
*/5 * * * * curl -o /workspace/category.json https://dnnmagento.addwebprojects.com/rest/V1/categorydata

# An empty line is required at the end of this file for a valid cron file.
